# deploy/systemd/panshi-studio-publish.timer
[Unit]
Description=Publish the next ready Panshi Studio trading-day edition every five minutes

[Timer]
OnCalendar=*-*-* *:00/5:00
Persistent=true
AccuracySec=30s
Unit=panshi-studio-publish.service

[Install]
WantedBy=timers.target
